# roles/ffe.unbound_exporter/files/unbound_exporter.service
[Unit]
Description=Unbound Exporter

[Service]
User=unbound
ExecStart=/usr/local/bin/unbound_exporter

[Install]
WantedBy=default.target
